# docs/deploy/systemd/resolvrr-dev.service.example
[Unit]
Description=Resolvrr dev server
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/resolvrr/resolvrr
Environment=NODE_ENV=development
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStartPre=/usr/bin/rm -rf /home/resolvrr/resolvrr/.next/dev/cache/turbopack
ExecStart=/usr/local/bin/npm run dev
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
